# test_watchdog.timer
[Unit]
Description=Run test monitor every minute

[Timer]
Unit=test_watchdog.service
OnCalendar=*-*-* *:*:00

[Install]
WantedBy=timers.target
