# test_watchdog.service
[Unit]
Description=Test process monitor

[Service]
Type=oneshot
ExecStart=/usr/local/bin/test_watchdog.sh

[Install]
WantedBy=multi-user.target
